# scripts/services/ods-display-config.service
[Unit]
Description=ODS Display Configuration (xrandr)
After=ods-player-ATLAS.service

[Service]
Type=oneshot
User=root
Environment=DISPLAY=:0
ExecStart=/usr/local/bin/ods-display-config.sh
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
