# scripts/services/ods-dpms-enforce.timer
[Unit]
Description=Enforce DPMS off every 5 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
