# scripts/services/ods-enrollment-retry.timer
[Unit]
Description=Retry ODS Cloud enrollment every 30 minutes until registered

[Timer]
OnBootSec=120
OnUnitActiveSec=1800

[Install]
WantedBy=timers.target
